# crontab
# Place in root crontab
@reboot /home/klack/auto-plot/tmux_startup.sh
